# sv/agent/run
#!/bin/bash

#sv start zabbix || exit 1

source /etc/envvars
env

MY_IP=`ip route get 1 | awk '{print $NF;exit}'`
sed -i -e "s|# ListenIP=.*|ListenIP=$MY_IP|" /usr/local/etc/zabbix_agentd.conf

if [ -n "$ZABBIX_SERVER" ]; then
  #ZABBIX_IP=`dig +short $ZABBIX_SERVER`
  sed -i -e "s|Server=.*|Server=$ZABBIX_SERVER,127.0.0.1|" -e "s|ServerActive=.*|ServerActive=$ZABBIX_SERVER|" /usr/local/etc/zabbix_agentd.conf
fi

if [ -n "$HOSTMETADATA" ]; then
  sed -i -e "s|HostMetadata=.*|HostMetadata=Linux Active $HOSTMETADATA|" /usr/local/etc/zabbix_agentd.conf
fi

exec 2>&1
chpst -u zabbix zabbix_agentd -f
